# DEBIAN/postinst
#!/bin/bash
set -e
systemctl daemon-reload
systemctl enable --now smart-updater.timer