# etc/systemd/system/smart-updater.service
[Unit]
Description=Smart background APT and Flatpak updates

[Service]
Type=oneshot
ExecStart=/usr/local/bin/smart-updater.sh